# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=service ssh start -D

[program:server]
command=/root/rccmd/rccmdServer/bin/linux/rccmdServer -test
